# xfce/skel/.config/autostart/compiz.desktop
[Desktop Entry]
Exec=xcompmgr -c -f -n
Name=compiz
Type=Application
Version=1.0
