# xfce/skel/.config/autostart/locker.desktop
[Desktop Entry]
Exec=xautolock -locker slock
Name=locker
Type=Application
Version=1.0
X-LXQt-X11-Only=true
