# xfce/skel/.config/autostart/synclient2.desktop
[Desktop Entry]
Exec=synclient TapButton2=3
Name=synclient2
Type=Application
Version=1.0
