# xfce/skel/.config/autostart/synclient1.desktop
[Desktop Entry]
Exec=synclient TapButton1=1
Name=synclient1
Type=Application
Version=1.0
