# xfce/skel/.config/autostart/synclient3.desktop
[Desktop Entry]
Exec=synclient TapButton3=2
Name=synclient3
Type=Application
Version=1.0
